# net-dns/dns2tcp/files/init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

#Config is in /etc/dns2tcpd.conf

depend() {
	need net
}

start() {
	
	ebegin "Starting dns2tcp"
	start-stop-daemon --start --quiet \
	--exec /usr/bin/dns2tcpd -- -f /etc/dns2tcpd.conf ${DNS2TCP_OPTS}
	eend ${?}
}

stop() {
	ebegin "Stopping dns2tcp"
	start-stop-daemon --stop --quiet --name dns2tcpd
	eend ${?}
}
