# www-apps/xelpaste/files/xelpaste.cron.daily
#!/bin/sh

xelpastectl cleanup_snippets > /dev/null
