# mail-filter/pfixtools/files/postlicyd-initd
#!/sbin/runscript

depend() {
    need net
    use dns
}

POSTLICYD=/usr/sbin/postlicyd
PIDFILE=/var/run/postlicyd/pid
CONF=/etc/pfixtools/postlicyd.conf

function do_checkconf() {
  $POSTLICYD -c "$CONF"
  return "$?"
}

start() {
    ebegin "Starting postlicyd..."
    do_checkconf || { eend 1 "Invalid configuration"; return 1; }
    start-stop-daemon --start --quiet --pidfile $PIDFILE --exec $POSTLICYD -- -p "$PIDFILE" "$CONF" || eerror "Failed"
    eend ${?}
}

stop() {
    ebegin "Stopping postlicyd..."
    start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 --pidfile $PIDFILE --name "postlicyd"
    eend ${?}
}

reload() {
    ebegin "Reloading postlicyd"
    do_checkconf || { eend 1 "Invalid configuration"; return 1; }
    start-stop-daemon --stop --signal 1 --quiet --pidfile $PIDFILE --name postlicyd
    eend ${?}
}
